# {{ cookiecutter.project_name }}/deploy/etc/s6-overlay/s6-rc.d/worker/run
#!/command/with-contenv execlineb

cd /app
foreground { python manage.py setup_periodic_tasks }
python manage.py qcluster