# {{ cookiecutter.project_name }}/deploy/etc/s6-overlay/s6-rc.d/server/run
#!/command/with-contenv execlineb

cd /app
foreground { python manage.py createsuperuser --noinput }
gunicorn config.wsgi:application --config="deploy/gunicorn.conf.py"
