# {{ cookiecutter.project_name }}/deploy/etc/s6-overlay/s6-rc.d/server/run
#!/command/with-contenv execlineb

cd /app
foreground { falco setup-admin }
gunicorn wsgi:application --config="deploy/gunicorn.conf.py"
